# scripts/avatar_worker.service.sample
[Unit]
Description=Avatar Worker Script
StartLimitIntervalSec=0

[Service]
Restart=always
RestartSec=1
ExecStart=/usr/bin/php /scripts/avatar_worker.php

[Install]
WantedBy=default.target

